# packaging/systemd/agentd.service
[Unit]
Description=MCP-native agent runtime (agentd)
Documentation=https://github.com/agentd-dev/source-code/tree/main/docs
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
EnvironmentFile=-/etc/default/agentd
ExecStart=/usr/bin/agentd $AGENT_ARGS
Restart=on-failure
RestartSec=5s

# Graceful shutdown. TimeoutStopSec MUST exceed agentd's --drain-timeout
# (default 25s) or systemd will SIGKILL mid-drain and in-flight work errors;
# a clean drain exits 0, an overrun is the OS-set 143 (rfcs/0011 §5).
KillSignal=SIGTERM
TimeoutStopSec=60s

# --- Hardening ---------------------------------------------------
# Dynamic UID/GID: no persistent user needed. Pair with
# StateDirectory / LogsDirectory below for write-access paths.
DynamicUser=yes

# StateDirectory → /var/lib/agentd ; LogsDirectory → /var/log/agentd.
# Both are created with the correct dynamic-user perms on first
# start and kept across restarts.
StateDirectory=agentd
StateDirectoryMode=0750
LogsDirectory=agentd
LogsDirectoryMode=0750

# Filesystem namespace: the process can only write to /var/lib/agentd,
# /var/log/agentd, and /tmp (PrivateTmp).
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ReadOnlyPaths=/etc/agentd

# Kernel + namespace isolation.
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectClock=yes
ProtectHostname=yes
RestrictNamespaces=yes
RestrictRealtime=yes
LockPersonality=yes

# Note: agentd forks + supervises subagent processes (PR_SET_CHILD_SUBREAPER,
# PR_SET_PDEATHSIG, process-group kills). The @system-service set + the
# restrictions below permit that. Do NOT add ProtectControlGroups=yes if you use
# --cgroup (active enforcement needs cgroup-fs writes).
MemoryDenyWriteExecute=yes

# Syscall allowlist. @system-service is the standard set.
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
SystemCallArchitectures=native

# Networking: AF_UNIX for the intel/served-MCP sockets; INET/INET6 for an
# https:// intelligence endpoint (tls build) and the /metrics+/healthz surface.
# No raw sockets, no netlink. Add AF_VSOCK only for the vsock build.
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

# Drop all capabilities. The runtime never needs any.
CapabilityBoundingSet=
AmbientCapabilities=
NoNewPrivileges=yes

# If the /metrics+/healthz surface must bind a port <1024, enable via drop-in:
#   systemctl edit agentd
# and add: CapabilityBoundingSet=CAP_NET_BIND_SERVICE
# + AmbientCapabilities=CAP_NET_BIND_SERVICE under [Service].

[Install]
WantedBy=multi-user.target
