# examples/systemd-agentd.service
# agent as a long-lived reactive daemon under systemd.
# Install at /etc/systemd/system/agent-triage.service, then:
#   systemctl daemon-reload && systemctl enable --now agent-triage
[Unit]
Description=agent ticket triage (reactive)
After=network.target

[Service]
Type=simple
Environment=AGENT_INTELLIGENCE=unix:/run/intelligence.sock
# Secrets via EnvironmentFile (mode 0600), never on the command line or a config
# file (RFC 0011 §3.2). The file sets AGENT_INTELLIGENCE_TOKEN=...
EnvironmentFile=/etc/agentd/triage.env
ExecStart=/usr/local/bin/agentd \
  --mode reactive \
  --instruction-file /etc/agentd/triage.txt \
  --model my-model \
  --mcp tickets=mcp-server-tickets \
  --subscribe tickets://queue/inbound \
  --drain-timeout 25s \
  --health-file /run/agent/health

# systemd's analogue of terminationGracePeriodSeconds — MUST exceed
# --drain-timeout (25s) or systemd SIGKILLs mid-drain.
TimeoutStopSec=30
KillSignal=SIGTERM
Restart=on-failure
RestartSec=5
RuntimeDirectory=agent

[Install]
WantedBy=multi-user.target
